# medical-django.service
[Unit]
Description=Medical Django Application
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/medical
Environment=PATH=/var/www/medical/venv/bin
ExecStart=/var/www/medical/venv/bin/python manage.py runserver 127.0.0.1:8000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
